# shared_scripts/docker_cleanup/docker_cleanup.cron
0 2 * * 0 /usr/local/bin/docker-cleanup.sh